# src/webserver/nginx-consultpl/tests/confighcl2/conf/supervisord/nginx.conf
[program:nginx]
command=nginx -c /conf-nginx/nginx.conf
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx-reload]
command=nginx -c /conf-nginx/nginx.conf -s reload
autostart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=0 ; no need to stay running

[program:nginx-config]
command=bash -c "while inotifywait -e create,delete,modify,move -q  /conf-nginx/; do supervisorctl -c  %(here)s/init.conf restart nginx-reload; done"
