# src/webserver/nginx-consultpl/tests/confighcl2/conf/supervisord/php-fpm.conf
[program:php5-fpm]
command=php5-fpm --fpm-config /etc2/php5/php-fpm.conf

[program:php5-fpm-config]
command=bash -c "while inotifywait -e create,delete,modify,move -q  /etc2/php5/; do supervisorctl -c  %(here)s/init.conf restart php5-fpm; done"