# src/webserver/nginx-consultpl/tests/confighcl2/conf/supervisord/testecho.conf
[program:echo]
command=echo %(here)s
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=0 ; no need to stay running

[program:env]
command=/usr/bin/env
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
startsecs=0  ; no need to stay running