# tests/env/fatubuntu-supervisord/data/supervisor/app1.conf
[program:app1]
command=/app/app1.sh
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
#[program:nginx-config]
#command=bash -c "while inotifywait -e create,delete,modify,move -q /data/config/; do supervisorctl restart nginx; done"
