# src/webserver/nginx-consultpl/tests/confighcl/conf/supervisord/consul-template.conf
[program:consul-template]
command=consul-template -config /etc2/consul-templates/config.hcl

[program:consul-template-config]
command=bash -c "while inotifywait -e create,delete,modify,move -q /etc2/consul-templates/; do supervisorctl -c  %(here)s/init.conf restart consul-template; done"
